# roles/hyprland/files/systemd/swww-randomize.timer
[Unit]
Description=Randomize Wallpapers using swww

[Timer]
OnStartupSec=0min
OnUnitActiveSec=20min

[Install]
WantedBy=timers.target
